# intro-ansible/07-web-app-deploys/roles/webserver/templates/supervisor_app.conf.j2
[program:{{ app_name }}]
command={{ venv_dir }}/bin/gunicorn app:app
directory={{ app_dir }}
user={{ deploy_user }}
autostart=true
autorestart=true
redirect_stderr=true
